# agent/install/macos/com.example.syshealth.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
<key>Label</key><string>com.example.syshealth</string>
<key>ProgramArguments</key>
<array><string>/usr/bin/python3</string><string>/opt/syshealth-agent/agent.py</string></array>
<key>RunAtLoad</key><true/>
<key>KeepAlive</key><true/>
<key>StandardErrorPath</key><string>/var/log/syshealth.err</string>
<key>StandardOutPath</key><string>/var/log/syshealth.out</string>
<key>EnvironmentVariables</key>
<dict><key>SYSHEALTH_CONFIG</key><string>/opt/syshealth-agent/config.json</string></dict>
</dict></plist>